# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:php5fpm]
command=php-fpm7.0 --fpm-config /etc/php/7.0/fpm/php-fpm.conf -F
